# config/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
user=root

[program:nodejs]
command=/usr/bin/node /var/www/sites/node-api.com/index.js
user=root